# etc/supervisor-conf.d/supervisord-php.conf
[program:php]
command=start-php.sh
autorestart=true
autostart=true
